# config/super_ertc.conf
[program:ertc]
command=/var/ertc/env/bin/python /var/ertc/ertc/ertc.py
user=lj
directory=/var/ertc/ertc 
startsecs=3 
environment=PYTHONPATH=".:/var/ertc/ertc"

