# config/super_ertc.conf
[program:ertc]
command=/var/ertc/env/bin/python /var/ertc/ertc/ertc.py
user=lj
directory=/var/ertc/ertc # Directory to cd before starting process
startsecs=3 # 1 is default VMs might be a bit slow

